# 2-settings-refine/systemd/refine-icons.service
[Unit]
Description=Re-bake theme-aware refined System Settings section icons

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cocovox-refine-icons
