# 7-notifications/systemd/swaync-apply-scheme.service
[Unit]
Description=Swap swaync CSS to match the active KDE light/dark color scheme

[Service]
Type=oneshot
ExecStart=%h/.local/bin/swaync-apply-scheme.sh
