# 7-notifications/systemd/swaync-apply-scheme.path
[Unit]
Description=Watch KDE color-scheme changes to repaint swaync (light/dark)

[Path]
# kdeglobals is rewritten whenever the color scheme / look-and-feel changes.
PathChanged=%h/.config/kdeglobals
Unit=swaync-apply-scheme.service

[Install]
WantedBy=default.target
